# scripts/home-hq-smart.timer.example
# Home HQ — run the SMART health collector daily.
# Install: copy to /etc/systemd/system/home-hq-smart.timer, then:
#   sudo systemctl daemon-reload && sudo systemctl enable --now home-hq-smart.timer
[Unit]
Description=Home HQ — daily SMART drive health collection

[Timer]
OnCalendar=daily
# Catch up on the next boot if the machine was off at the scheduled time.
# (systemd has no trailing inline comments — keep them on their own line.)
Persistent=true

[Install]
WantedBy=timers.target
