# scripts/home-hq-tailscale.timer.example
# Home HQ — refresh Tailscale mesh status periodically.
# Install: copy to /etc/systemd/system/home-hq-tailscale.timer, then:
#   sudo systemctl daemon-reload && sudo systemctl enable --now home-hq-tailscale.timer
# NOTE: systemd does not support trailing inline comments on a directive line —
# keep comments on their own lines, or the value parse fails.
[Unit]
Description=Home HQ — periodic Tailscale status refresh

[Timer]
# Every 5 minutes (wall clock) — keeps online/offline states reasonably fresh.
OnCalendar=*:0/5
# Catch up on the next boot if the machine was off at the scheduled time.
Persistent=true

[Install]
WantedBy=timers.target
